# blueprints/secops-ai-migration-helper/recommender_curated_community/resources/example_recommendation_curated_community_rules.json
[
  {
    "ucid": "9901",
    "title": "Rare Scripting Software Detected",
    "description": "Because powershell and cmd prompt are so closely watched and have so many alarms put up to catch misuse, hackers are starting to avoid using them. This detects scripting tools that are utilised and found in the Event log that aren't approved software by your environment.",
    "curated rules": "Applied Threat Intelligence - Curated Prioritization/High Priority Host Indicators/ATI High Priority Rule Match for File IoCs (principal.process.file.sha256),Applied Threat Intelligence - Curated Prioritization/High Priority Host Indicators/ATI High Priority Rule Match for File IoCs (target.process.file.sha256),Applied Threat Intelligence - Curated Prioritization/High Priority Host Indicators/ATI High Priority Rule Match for File IoCs (target.file.sha256),Applied Threat Intelligence - Curated Prioritization/High Priority Host Indicators/ATI High Priority Rule Match for File IoCs (about.file.sha256)",
    "curated rules coverage": "very good",
    "curated rationale": "The security objective is to detect malicious scripting activity that avoids common tools. The suggested curated rules are part of the 'Applied Threat Intelligence - Curated Prioritization' set and will trigger if the hash of a running process or a file on disk matches a known malicious indicator from Google's threat intelligence. This effectively detects known malicious or rare scripting tools without relying on potentially noisy keyword matching. Log Data Sources: EDR logs such as Microsoft Sysmon (WINDOWS_SYSMON), Crowdstrike Falcon (CS_EDR), SentinelOne (SENTINEL_EDR), Carbon Black (CB_EDR).",
    "community rules": "recon_suspicious_commands_cisa_report,mitre_attack_T1140_encoded_powershell_command,gcti_remote_access_tools,recon_environment_enumeration_system_cisa_report,google_safebrowsing_file_process_creation",
    "community rules coverage": "very good",
    "community rationale": "The selected community rules collectively detect the execution of suspicious or unauthorized tools and scripts, aligning with the input rule's objective. The `recon_suspicious_commands_cisa_report` and `recon_environment_enumeration_system_cisa_report` rules use specific, intelligence-driven command-line patterns to identify a wide range of malicious activities. The `mitre_attack_T1140_encoded_powershell_command` rule targets a common evasion technique used with scripting tools. Furthermore, `gcti_remote_access_tools` and `google_safebrowsing_file_process_creation` leverage threat intelligence feeds to detect known malicious tools by their file hashes, offering a more robust detection method than the input rule's generic keyword matching. Log Data Sources: Microsoft Sysmon, Microsoft Windows Events, CrowdStrike, SentinelOne, Microsoft Defender ATP."
  },
  {
    "ucid": "9902",
    "title": "Rare or Never Seen Events EventID Logged",
    "description": "This requires a watchlist of every event ID from this table and one of  rare event IDs",
    "curated rules": "Risk Analytics for UEBA/New Login by User to Device/First Time Interactive Service Account Login by Principal User ID and Target Service Account",
    "curated rules coverage": "partially",
    "curated rationale": "The curated rule detects the first time a user interactively logs in with a service account, which aligns with the input rule's objective of identifying new or rare events. While the input rule is generic and looks for any new EventID, the curated rule applies the 'first seen' concept to a specific, high-value security scenario. Log Data Sources: Windows Security Events (WINDOWS_EVENT_LOG), Linux logs (NIX_SYSTEM), or various cloud/SaaS application logs that produce USER_LOGIN events.",
    "community rules": "hash_prevalence,domain_prevalence,ip_target_prevalence,network_connection_first_seen_in_past_day,low_prevalence_hash_on_process_launch_low_prevalence_domain_accessed",
    "community rules coverage": "very good",
    "community rationale": "The input rule's objective is to detect rare or never-before-seen Windows Event IDs, which is a form of anomaly detection based on event frequency. The suggested community rules achieve the same security objective by applying prevalence and first-seen logic to other critical security data points. 'hash_prevalence', 'domain_prevalence', and 'ip_target_prevalence' directly use prevalence calculations to find rare file hashes, domains, and IP addresses, respectively. 'network_connection_first_seen_in_past_day' detects new IP addresses, aligning with the 'never seen' aspect of the input rule. 'low_prevalence_hash_on_process_launch_low_prevalence_domain_accessed' provides a high-fidelity signal by combining two prevalence checks. Collectively, these rules provide robust coverage for detecting rare and potentially malicious activities through anomaly detection. Log Data Sources: EDR logs (e.g., Microsoft Sysmon, CrowdStrike), DNS logs, and Network/Firewall logs (e.g., Zeek, GCP Firewall)."
  },
  {
    "ucid": "9903",
    "title": "New Run/RunOnce Scheduled Task",
    "description": "Shell-Core EventID 9706/9707/9708",
    "curated rules": "N/A",
    "curated rules coverage": "no",
    "curated rationale": "No equivalent curated rule was found that specifically detects persistence through the creation or modification of Run/RunOnce registry keys. The provided curated rules focus on IOC matching or high-level UEBA detections like anomalous logins and network traffic, which do not cover this specific persistence technique.",
    "community rules": "mitre_attack_T1053_005_windows_creation_of_scheduled_task,suspicious_unusual_location_lnk_file",
    "community rules coverage": "partially",
    "community rationale": "The input rule detects persistence via Run/RunOnce registry keys (T1547.001). The community rules cover similar persistence objectives. 'mitre_attack_T1053_005_windows_creation_of_scheduled_task' detects the creation of scheduled tasks (T1053.005), another common method for autostart execution. 'suspicious_unusual_location_lnk_file' detects the creation of shortcut files in startup folders, which is also part of the T1547 technique family. While the specific implementation differs, the overall goal of detecting persistence is addressed. Log Data Sources: Microsoft Sysmon, SentinelOne, or other EDR logs providing process launch and file events."
  },
  {
    "ucid": "9904",
    "title": "Unusual binary shutting down endpoint",
    "description": "Malware may cause shutdowns on this EID as well as EventID 6006 too",
    "curated rules": "Applied Threat Intelligence - Curated Prioritization/Active Breach Priority Host Indicators/ATI Active Breach Rule Match for File IoCs (principal.process.file.sha256),Applied Threat Intelligence - Curated Prioritization/High Priority Host Indicators/ATI High Priority Rule Match for File IoCs (target.file.sha256),Applied Threat Intelligence - Curated Prioritization/Medium Priority Host Indicators/ATI Medium Priority Rule Match for File IoCs (principal.process.file.sha256),Applied Threat Intelligence - Curated Prioritization/Medium Priority Host Indicators/ATI Medium Priority Rule Match for File IoCs (target.file.sha256)",
    "curated rules coverage": "very good",
    "curated rationale": "The security objective is to detect malware initiating a system shutdown. While there isn't a specific curated rule for the shutdown event (EventID 1074), the Applied Threat Intelligence rules provide excellent coverage by detecting the malicious binary itself based on its file hash when it is executed or written to disk. This often provides an earlier detection of the threat. Log Data Sources: EDR (Carbon Black, Crowdstrike, SentinelOne, Microsoft Sysmon).",
    "community rules": "N/A",
    "community rules coverage": "no",
    "community rationale": "N/A"
  }
]